# build/sysroot/usr/lib/systemd/system/casaos-local-storage.service
[Unit]
After=casaos-gateway.service
After=casaos-message-bus.service
After=casaos-user-service.service
After=casaos.service
Before=docker.service
ConditionFileNotEmpty=/etc/casaos/local-storage.conf
Description=CasaOS Local Storage Service

[Service]
ExecStartPre=/usr/bin/casaos-local-storage -v
ExecStart=/usr/bin/casaos-local-storage -c /etc/casaos/local-storage.conf
PIDFile=/var/run/casaos/local-storage.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
